# BackEnd/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/bin/sh -c "php artisan serve --host=0.0.0.0"

[program:web-socket]
command=/bin/sh -c "php artisan websocket:serve"

